# claude/systemd/git-hygiene.service
[Unit]
Description=Daily branch-hygiene check + weekly auto-clean (jckeen/dotfiles)
Documentation=https://github.com/jckeen/dotfiles
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=%h/dev/dotfiles/claude/scripts/hygiene-cron.sh

# ─── Hardening ───────────────────────────────────────────────
# systemd directives do NOT accept inline `#` comments on the same line;
# every annotation below must live on its own line.

# Block setuid/setgid escalation.
NoNewPrivileges=true
# Private /tmp and /var/tmp for this unit.
PrivateTmp=true
# Entire FS becomes read-only except API filesystems (/proc, /sys) and ReadWritePaths.
ProtectSystem=strict
# NOTE: ProtectHome=tmpfs would break this service — it reads scripts from $HOME/dev/dotfiles
# and writes state under $HOME/.local/state/hygiene/. Using read-only + explicit ReadWritePaths.
ProtectHome=read-only
# hygiene-cron writes status.json + cron.log under this dir; needs to be writable.
ReadWritePaths=%h/.local/state/hygiene
# TODO: verify ReadWritePaths covers anything git-hygiene clean touches outside ~/dev (git refs, locks, etc.)
# Drop AF_NETLINK, AF_PACKET, raw sockets — keep only UNIX + IPv4/IPv6 (gh API + git fetch).
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# Forbid creating new user/mount/net/uts/ipc/pid/cgroup namespaces.
RestrictNamespaces=true
# Lock personality(2) — block legacy ABI / x86_32 emulation switching.
LockPersonality=true
# Block W+X memory — bash + gh + git do not JIT, so this is safe here.
MemoryDenyWriteExecute=true
# Whitelist common syscalls for typical services; everything else returns EPERM.
SystemCallFilter=@system-service
# Return EPERM for filtered syscalls (friendlier than SIGSYS for app-level handling).
SystemCallErrorNumber=EPERM
